# crontab
*/1 * * * * /usr/local/bin/python /src/data/data_gen.py
